# deploy/lodestar-backup.service
# systemd oneshot — snapshot every tenant's claims.log (run on a schedule by
# lodestar-backup.timer; this unit does the work).
#   sudo cp deploy/lodestar-backup.service /etc/systemd/system/
#   sudo cp deploy/lodestar-backup.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lodestar-backup.timer
[Unit]
Description=Lodestar per-tenant claims.log backup

[Service]
Type=oneshot
User=lodestar
Group=lodestar
Environment=LODESTAR_TENANT_ROOT=/var/lib/lodestar/tenants
Environment=LODESTAR_BACKUP_DIR=/var/lib/lodestar/backups
ExecStart=/opt/lodestar/deploy/backup.sh

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/lodestar
